# ubuntu-health-monitor.service
[Unit]
Description=Ubuntu Health Monitor Flask API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/path/to/Ubuntu-Health-Monitor
ExecStart=/path/to/Ubuntu-Health-Monitor/venv/bin/gunicorn -b 0.0.0.0:5000 app:app
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ubuntu-health-monitor
Environment="PATH=/path/to/Ubuntu-Health-Monitor/venv/bin"

[Install]
WantedBy=multi-user.target
